# self-healing-docker-nginx-setup/self_heal.cron
* * * * * root /usr/local/bin/self_heal.sh >> /var/log/self_heal.log 2>&1
